# chmod-itunes/ch.wit.fix-iTunes-permissions.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ch.wit.fix-iTunes-permissions</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Program</key>
    <string>/usr/local/sbin/fix-iTunes-permissions.sh</string>

    <key>StartInterval</key>
    <integer>127</integer>
</dict>
</plist>
